# pkgbuild/alfred.service
[Unit]
Description=Alfred
After=network.target

[Service]
ExecStart=/usr/local/bin/alfred
Restart=always

[Install]
WantedBy=multi-user.target
