# rpmbuild/SOURCES/alfred.service
[Unit]
Description=Alfred
After=network.target

[Service]
ExecStart=/usr/bin/alfred
Restart=always

[Install]
WantedBy=multi-user.target
