# rpmbuild/SOURCES/botman.service
[Unit]
Description=Botman
After=network.target

[Service]
ExecStart=/usr/bin/botman
Restart=always

[Install]
WantedBy=multi-user.target
